# cron.d/crontab
* * * * * php /var/www/artisan schedule:run